# system_files/usr/lib/systemd/system/armada-waydroid-input.path
[Unit]
Description=Watch for Waydroid Android session startup

[Path]
PathChanged=/run/waydroid-lxc/network_up
Unit=armada-waydroid-input.service

[Install]
WantedBy=multi-user.target
